# misc-home/home/.config/systemd/user/backup-webdav.service
[Unit]
Description=Backup of /tmp/webdav

[Service]
Type=simple
ExecStart=/usr/local/bin/backup-webdav.sh
