# misc-home/usr/local/lib/systemd/system/renew-certs.service
[Unit]
Description=Renew my certbot certs

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
